# others/governance/crontab
#######################
## GOVERNANCE ##
#######################
0,15,30,45 * * * * /home/cob/bin/sendMsg2IM clock clockTick governance > /dev/null