# others/customize.governance/crontab
#######################
## GOVERNANCE ##
#######################
0,15,30,45 * * * * /usr/local/bin/sendMsg2IM clock clockTick governance > /dev/null